# host/boxcutter-host.service
[Unit]
Description=Boxcutter Host Control Plane
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/boxcutter-host run
Environment=BOXCUTTER_REPO=/home/budda/boxcutter
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
