# host/boxcutter-host.service
[Unit]
Description=Boxcutter Host Control Plane
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/boxcutter-host run
# Set BOXCUTTER_REPO to the repo checkout for dev mode.
# In prod mode, the binary uses /var/lib/boxcutter as the data directory.
#Environment=BOXCUTTER_REPO=/path/to/boxcutter
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
