# systemd/boxcutter.service
[Unit]
Description=Boxcutter (SlicerVM daemon)
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/slicer up /etc/boxcutter/agent-dev.yaml
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
