# systemd/boxcutter-proxy-sync.service
[Unit]
Description=Boxcutter proxy sync (VM services → Caddy)
After=network.target caddy.service

[Service]
Type=simple
ExecStart=/usr/local/bin/boxcutter-proxy-sync
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
