# systemd/boxcutter-gateway.service
[Unit]
Description=Boxcutter gateway (noop — gateway runs via SSH ForceCommand)
After=network.target

[Service]
Type=oneshot
ExecStart=/bin/true
RemainAfterExit=yes
