# web/systemd/boxcutter-tailscale-cert.service
[Unit]
Description=Rotate Tailscale TLS certificate for nginx

[Service]
Type=oneshot
ExecStart=/usr/local/bin/boxcutter-tailscale-cert
